# oalworkers.conf
[program:oal_worker]
command=python /app/worker.py
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
numprocs=1
redirect_stderr=true
stopwaitsecs=60
